# scripts/systemd/auto-poly-btc-monitor-freshness.service
[Unit]
Description=Auto Polymarket - BTC 1s Freshness Watchdog
After=auto-poly-btc-monitor.service

[Service]
Type=oneshot
WorkingDirectory=/root/auto_polymarket
Environment=BTC_1S_MAX_STALE_SEC=20
Environment=BTC_1S_RECHECK_SEC=45
ExecStart=/bin/bash scripts/check_btc_1s_freshness.sh
