# scripts/systemd/auto-poly-advisory-calibration.timer
[Unit]
Description=Auto Polymarket - Advisory Calibration Timer (every 6h)

[Timer]
OnBootSec=10min
OnUnitActiveSec=6h
AccuracySec=2min
Unit=auto-poly-advisory-calibration.service

[Install]
WantedBy=timers.target
